# agent/system/f2b-agent.service
[Unit]
Description=Fail2BanEntreprise Agent — Ban sync
Documentation=https://f2b.callcenter-erp.com/servers
After=network-online.target fail2ban.service
Wants=network-online.target fail2ban.service

[Service]
Type=simple
ExecStart=/usr/local/bin/f2b-agent sync-loop
Restart=on-failure
RestartSec=15

# Suppress stdout from reaching syslog/journal — agent writes its own log file.
# This prevents log duplication into /var/log/syslog and /var/log/daemon.log.
StandardOutput=null
StandardError=journal

# Harden the service
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
